# luci-app-live/root/etc/uci-defaults/node
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@live[-1]
	add ucitrack live
	set ucitrack.@live[-1].init=live
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
